# rc.local
#!/bin/bash

#if rc.local doesnt exist, just creating it works in ubuntu - with a chmod
#cp ./rc.local /etc/rc.local

#create a tmux session as myself
sudo -i -H -u george tmux new-session -d -s boot-session

exit 0
